# exploresc-server/tileserver/supervisorprograms.ini
[program:postgresql]
command=runuser -l postgres -c "exec /usr/lib/postgresql/9.5/bin/postgres -D /etc/postgresql/9.5/main"
autostart=true
autorestart=unexpected

[program:apache]
command=apachectl -D FOREGROUND
autostart=true
autorestart=true

[program:renderd]
command=renderd -f -c /usr/local/etc/renderd.conf
autostart=true
autorestart=true

